# data/motd/fwupd-refresh.service.in
[Unit]
Description=Refresh fwupd metadata
Documentation=man:fwupdmgr(1)
After=network.target network-online.target systemd-networkd.service NetworkManager.service connman.service

[Service]
Type=oneshot
StandardOutput=file:@motd_file@
ExecStart=@bindir@/fwupdmgr refresh
ExecStart=@bindir@/fwupdmgr get-updates --no-metadata-check --no-unreported-check --quiet
