# data/motd/fwupd-refresh.service.in
[Unit]
Description=Refresh fwupd metadata and update motd
Documentation=man:fwupdmgr(1)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
CacheDirectory=fwupdmgr
StandardError=null
@user@
RestrictAddressFamilies=AF_NETLINK AF_UNIX AF_INET AF_INET6
SystemCallFilter=~@mount
ProtectKernelModules=yes
ProtectControlGroups=yes
SuccessExitStatus=2 101
ExecStart=@bindir@/fwupdtool refresh --json
KeyringMode=private
LockPersonality=yes
MemoryDenyWriteExecute=yes
NoNewPrivileges=no
PrivateDevices=no
PrivateTmp=true
ProtectClock=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelTunables=yes
ProtectProc=invisible
ProtectSystem=full
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
Environment="GLIBC_TUNABLES=glibc.cpu.hwcaps=SHSTK"
RestrictAddressFamilies=AF_NETLINK AF_UNIX AF_INET AF_INET6
