# data/motd/fwupd-refresh.service.in
[Unit]
Description=Refresh fwupd metadata and update motd
Documentation=man:fwupdmgr(1)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
CacheDirectory=fwupdmgr
StandardError=null
@user@
RestrictAddressFamilies=AF_NETLINK AF_UNIX AF_INET AF_INET6
SystemCallFilter=~@mount
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
SuccessExitStatus=2
ExecStart=@bindir@/fwupdmgr refresh
