# data/motd/fwupd-refresh.service.in
[Unit]
Description=Refresh fwupd metadata and update motd
Documentation=man:fwupdmgr(1)
After=network.target network-online.target systemd-networkd.service NetworkManager.service connman.service

[Service]
Type=oneshot
RuntimeDirectory=@motd_dir@
CacheDirectory=fwupd
RuntimeDirectoryPreserve=yes
StandardError=null
ExecStart=@bindir@/fwupdmgr refresh
ExecStart=@bindir@/fwupdmgr get-updates --log @motd_file@
DynamicUser=yes
RestrictAddressFamilies=AF_NETLINK AF_UNIX AF_INET AF_INET6
SystemCallFilter=~@mount
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
