# data/motd/fwupd-refresh.timer
[Unit]
Description=Refresh fwupd metadata regularly
ConditionVirtualization=!container

[Timer]
OnCalendar=*-*-* *:00:00
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
