# data/fwsignd.service.in
[Unit]
Description=Firmware signing daemon

[Service]
Type=simple
# you can set the user to non root if you want, but we're not then able to set
# the owner to the correct value on files written to the destination directory
#User=fwsignd
ExecStart=@servicedir@/fwupd/fwsignd
PrivateNetwork=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
