# data/fwupd.service.in
[Unit]
Description=Firmware update daemon
Documentation=http://www.fwupd.org/
After=dbus.service
Before=gdm.service

[Service]
Type=dbus
BusName=org.freedesktop.fwupd
ExecStart=@servicedir@/fwupd/fwupd
MemoryDenyWriteExecute=yes
PrivateTmp=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectKernelModules=yes
ProtectSystem=full
RestrictAddressFamilies=AF_NETLINK AF_UNIX
RestrictRealtime=yes
ReadWritePaths=@localstatedir@/lib/fwupd -@bootdir@
