# data/fwupd-offline-update.service.in
[Unit]
Description=Updates device firmware whilst offline
Documentation=man:fwupdmgr
ConditionPathExists=@localstatedir@/lib/fwupd/pending.db
DefaultDependencies=false
Requires=sysinit.target dbus.socket
After=sysinit.target system-update-pre.target dbus.socket systemd-journald.socket
Before=shutdown.target system-update.target

[Service]
Type=oneshot
ExecStart=@libexecdir@/fwupd/fwupdoffline
FailureAction=reboot
