# data/fwupdagent.timer
[Unit]
Description=Run fwupdagent daily and on boot

[Timer]
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
