# data/fwupd-offline-update.service.in
[Unit]
Description=Updates device firmware whilst offline
OnFailure=reboot.target

[Service]
ExecStart=@bindir@/fwupdmgr update-prepared
