# data/fwupd.service.in
[Unit]
Description=Firmware update daemon
Documentation=http://www.fwupd.org/
After=dbus.service
Before=gdm.service
ConditionPathExists=/var/lib/fwupd/pending.db

[Service]
Type=dbus
BusName=org.freedesktop.fwupd
ExecStart=@servicedir@/fwupd/fwupd
PrivateTmp=yes
