# data/fwupd-offline-update.service.in
[Unit]
Description=Updates device firmware whilst offline
Documentation=man:fwupdmgr
OnFailure=reboot.target

[Service]
ExecStart=@bindir@/fwupdmgr install-prepared
