# data/fwupd.service.in
[Unit]
Description=Firmware update daemon
After=dbus.service
Before=gdm.service
ConditionPathExists=/var/lib/fwupd/pending.db

[Service]
Type=dbus
BusName=org.freedesktop.fwupd
ExecStart=@servicedir@/fwupd/fwupd
PrivateNetwork=yes
PrivateTmp=yes
