# data/fwupd.service.in
[Unit]
Description=Firmware update daemon
Documentation=http://www.fwupd.org/
After=dbus.service
Before=gdm.service
ConditionPathExists=/var/lib/fwupd/pending.db

[Service]
Type=dbus
BusName=org.freedesktop.fwupd
ExecStart=@servicedir@/fwupd/fwupd
MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
PrivateTmp=yes
PrivateUsers=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=full
RestrictAddressFamilies=AF_NETLINK AF_UNIX
RestrictRealtime=yes
