# data/fwupd-offline-update.service.in
[Unit]
Description=Updates device firmware whilst offline
Documentation=man:fwupdmgr
OnFailure=reboot.target
ConditionPathExists=/var/lib/fwupd/pending.db

[Service]
ExecStart=@bindir@/fwupdmgr install-prepared
