# data/org.freedesktop.fwupd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>org.freedesktop.fwupd</string>
  <key>ProgramArguments</key>
  <array>
    <string>@libexecdir@/fwupd/fwupd</string>
  </array>
  <key>Sockets</key>
  <dict>
    <key>Listener</key>
    <dict>
      <key>SockPathName</key>
      <string>@dbus_socket_address@</string>
      <key>SockPathMode</key>
      <integer>438</integer><!-- 0666 -->
    </dict>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
</dict>
</plist>
